# scripts/taskenize.conf
description "node.js taskenize server"
author      "[yg] - http://gorvitovskiy.com"

# Used to Be: Start on Startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/home/ubuntu/"

    exec sudo -u ubuntu /usr/bin/node /home/ubuntu/Taskenize/WebServer/bin/app.js >> /var/log/taskenize.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script
